# com.joshuaterminal.app.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

  <key>Label</key>
  <string>com.joshuaterminal.app</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/<YourUserName>/.local/bin/uv</string>
    <string>run</string>
    <string>--project</string>
    <string>/Users/<YourUserName>/JoshuaTerminal</string>
    <string>python</string>
    <string>/Users/<YourUserName>/JoshuaTerminal/app.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/<YourUserName>/JoshuaTerminal</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>5</integer>

  <key>StandardOutPath</key>
  <string>/Users/<YourUserName>/Library/Logs/JoshuaTerminal/stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/<YourUserName>/Library/Logs/JoshuaTerminal/stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/<YourUserName></string>
  </dict>

</dict>
</plist>
